# windrose-monitor.service
[Unit]
Description=Windrose Server Monitor
Documentation=https://github.com/VegyBS/windrose-monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=windrose-monitor
Group=windrose-monitor

# Ensure predictable working directory
WorkingDirectory=/var/lib/windrose-monitor

# Load environment variables
EnvironmentFile=/etc/windrose-monitor/.env

# Use the venv Python selected during installation
ExecStart=/var/lib/windrose-monitor/venv/bin/python /usr/local/bin/windrose-monitor

# Restart behaviour
Restart=on-failure
RestartSec=10
SyslogIdentifier=windrose-monitor

# --- File Logging ---
StandardOutput=append:/var/log/windrose-monitor/windrose-monitor.log
StandardError=append:/var/log/windrose-monitor/windrose-monitor.log

# --- Security Hardening ---
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectClock=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
